# deploy/com.mitchens.omniminer-reports-rebuild.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  OPTIONAL ALTERNATIVE — NOT the active mechanism.
  The reports-site publish is bundled into the daily-maintenance job (daily 08:00,
  com.mitchens.daily-maintenance → rebuild_and_push.sh) as of 260606. That gives a
  daily rebuild with no extra daemon. Load THIS standalone watcher ONLY if you want
  near-instant publishing (15-min poll + on-change) instead of once-daily — e.g. so the
  Telegram report deep-link resolves within minutes rather than by the next morning.
  If you load this, you may want to drop the call from daily-maintenance.sh to avoid
  two triggers (harmless if both run — the mtime guard makes the loser a no-op).

  Install (only if choosing near-instant over the daily bundle):
    cp deploy/com.mitchens.omniminer-reports-rebuild.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.mitchens.omniminer-reports-rebuild.plist
  Disable:
    launchctl unload ~/Library/LaunchAgents/com.mitchens.omniminer-reports-rebuild.plist

  WatchPaths fires on filesystem events (may be unreliable on the GDrive FUSE mount);
  StartInterval is the reliable backstop. The script itself no-ops cheaply when no
  source .md is newer than the last build, so a short interval is inexpensive.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mitchens.omniminer-reports-rebuild</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/mitchens/Local/APPS/omniminer-reports/rebuild_and_push.sh</string>
    </array>
    <key>WatchPaths</key>
    <array>
        <string>/Users/mitchens/Library/CloudStorage/GoogleDrive-henspham@gmail.com/My Drive/_SYNC/OMNIMINER</string>
    </array>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/mitchens/Local/APPS/omniminer-reports/rebuild.launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mitchens/Local/APPS/omniminer-reports/rebuild.launchd.log</string>
</dict>
</plist>
